# files/systemd/system/dovecot.service.d/environment.conf
[Service]
EnvironmentFile=/etc/default/dovecot
